# tournament-manager.service.example
[Unit]
Description=Valorant Tournament Manager Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/your/valorant-tournament
ExecStart=/path/to/your/valorant-tournament/venv/bin/python /path/to/your/valorant-tournament/main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
